# scripts/nix-guard.service
[Unit]
Description=Nix Build CPU Guard
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/raed/projects/phonebridge/scripts/nix-guard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
